# tbxpost.service
[Unit]
Description=TBXpost — local Telegram bot daemon
Documentation=https://github.com/davesovm/TBXpost
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sovm
WorkingDirectory=/home/sovm/Desktop/sovm%27s/TBXpost
ExecStart=/home/sovm/Desktop/sovm%27s/TBXpost/venv/bin/python /home/sovm/Desktop/sovm%27s/TBXpost/bot/run_bot.py
Restart=on-failure
RestartSec=15s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tbxpost
EnvironmentFile=/home/sovm/Desktop/sovm%27s/TBXpost/.env

[Install]
WantedBy=default.target
